# systemd/bt-wake-sync.service
[Unit]
Description=Sync native Bluetooth wake peers to Pico wake controller
After=bluetooth.service
Requires=bluetooth.service

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/bt-wake-sync.py --serial /dev/serial/by-id/usb-Codex_Bluetooth_USB_Wake_Device_0009-if00 --adapter 0
